# roles/zot/templates/zot.service.j2
[Unit]
Description=Zot OCI Registry
Documentation=https://zotregistry.io
After=network.target

[Service]
Type=simple
User=zot
Group=zot
ExecStart={{ zot_binary_dest }}/{{ zot_binary_final_name }} serve {{ zot_config_dir }}/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zot

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ zot_data_dir }} {{ zot_home_dir }}
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target